# data/intel_lpmd.service.in
[Unit]
Description= Intel Linux Energy Optimizer (lpmd) Service
ConditionVirtualization=no
StartLimitInterval=200
StartLimitBurst=5

[Service]
Type=dbus
SuccessExitStatus=2
BusName=org.freedesktop.intel_lpmd
ExecStart=@sbindir@/intel_lpmd --systemd --dbus-enable
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
Alias=org.freedesktop.intel_lpmd.service
