# data/intel_lpmd.service.in
[Unit]
Description= Intel Linux Energy Optimizer (lpmd) Service
ConditionVirtualization=no

[Service]
Type=simple
SuccessExitStatus=2
BusName=org.freedesktop.intel_lpmd
ExecStart=/usr/bin/intel_lpmd --systemd --dbus-enable
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=org.freedesktop.intel_lpmd.service
